# syncpod.service
[Unit]
Description=syncpod

[Service]
Type=simple
User=derrick
ExecStart=/home/derrick/.bin/syncpod
